# services/eink.service
[Unit]
Description=EINK Display Service

After=syslog.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/eink.py
Restart=always
RestartSec=5
KillMode=process
User=pi

[Install]
WantedBy=multi-user.target
